# sunreye/rootfs/etc/s6-overlay/s6-rc.d/init-db-ready/run
#!/usr/bin/with-contenv bashio
# Gate everything DB-dependent on the database actually answering — works for
# both the embedded postgres (just started) and an external one (may still be
# booting when this addon starts).
set -euo pipefail

for _ in $(seq 1 60); do
    if pg_isready -d "${DATABASE_URL}" -q 2>/dev/null; then
        # Embedded mode: make sure the application database exists (initdb
        # only creates the postgres/template databases).
        if [ "${DB_EMBEDDED}" = "true" ]; then
            s6-setuidgid postgres psql -X -tAc "SELECT 1 FROM pg_database WHERE datname = 'sunreye'" | grep -q 1 \
                || s6-setuidgid postgres createdb sunreye
        fi
        exit 0
    fi
    sleep 1
done

bashio::log.fatal "Database at DATABASE_URL did not become ready within 60s"
exit 1
