# sunreye/rootfs/etc/s6-overlay/s6-rc.d/svc-server/run
#!/usr/bin/with-contenv bashio
# SunReye core engine (compiled bun binary) on 127.0.0.1:3000.
export PORT=3000
exec s6-setuidgid sunreye /opt/sunreye/server
