# sunreye/rootfs/etc/s6-overlay/s6-rc.d/init-postgres/run
#!/usr/bin/with-contenv bashio
# Prepare the embedded PostgreSQL datadir under /data (persisted by the
# Supervisor across upgrades and included in HA backups). Skipped entirely
# when an external database is configured.
set -euo pipefail

if [ "${DB_EMBEDDED}" != "true" ]; then
    bashio::log.info "External database configured — embedded PostgreSQL disabled"
    exit 0
fi

mkdir -p /data/backups
install -d -o postgres -g postgres -m 700 /data/postgres
install -d -o postgres -g postgres /var/run/postgresql

if [ ! -s /data/postgres/PG_VERSION ]; then
    bashio::log.info "Initializing PostgreSQL ${PG_MAJOR} data directory (first start)"
    s6-setuidgid postgres initdb -D /data/postgres -U postgres \
        --encoding=UTF8 --auth-local=peer --auth-host=trust
fi

# Managed settings live in their own include, rewritten on every start so
# upgrades can retune existing data directories (settings appended by earlier
# addon versions are overridden — the include comes last, and last wins).
cat > /data/postgres/sunreye.conf <<'EOF'
listen_addresses = '127.0.0.1'
shared_preload_libraries = 'timescaledb'
timescaledb.telemetry_level = off
# Headroom for the TimescaleDB job scheduler (cagg refresh x3, compression,
# retention): the stock max_worker_processes=8 makes job launches fail with
# "failed to start a background worker".
max_worker_processes = 16
timescaledb.max_background_workers = 8
# The image ships without LLVM (JIT machinery pruned for size); this must
# stay off or planner attempts to JIT would fail to load llvmjit.so.
jit = off
EOF
chown postgres:postgres /data/postgres/sunreye.conf
grep -q "include = 'sunreye.conf'" /data/postgres/postgresql.conf \
    || echo "include = 'sunreye.conf'" >> /data/postgres/postgresql.conf

# Base-image updates can shift uids; make sure postgres still owns its datadir.
if [ "$(stat -c %U /data/postgres)" != "postgres" ]; then
    chown -R postgres:postgres /data/postgres
fi

# Hard guard: a datadir is only readable by the pg major that created it.
# Refuse to boot instead of letting postgres scribble error loops — the addon
# shows up as stopped with this message as the last log lines.
datadir_major="$(cat /data/postgres/PG_VERSION)"
if [ "${datadir_major}" != "${PG_MAJOR}" ]; then
    bashio::log.fatal "Data directory is PostgreSQL ${datadir_major}, this addon ships PostgreSQL ${PG_MAJOR}."
    bashio::log.fatal "A migration between PostgreSQL majors requires a dedicated transition release —"
    bashio::log.fatal "do not downgrade the addon; check the SunReye release notes for the upgrade path."
    exit 1
fi
